# vc3-omnibus-website/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:vc3-infoservice]
command=/usr/bin/vc3-infoservice --runas=vc3 --conf=/etc/vc3/vc3-infoservice.conf --log=/var/log/vc3/infoservice.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:vc3-master]
command=/usr/bin/vc3-master --runas=vc3 --conf=/etc/vc3/vc3-master.conf --log=/var/log/vc3/master.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi.d/vc3.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
